# crontab
0 0,6,12,18 * * * /bin/sh /backup.sh >> /var/log/cron.log
0 3,9,15,21 * * * /bin/sh /sync.sh >> /var/log/cron.log
# this line is supposed to be here
